# crontab
0 21 * * * python /script.py >> /logfile